# services/node.service
[Unit]
Description=Node Broadcast - sensor readings via LoRa
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 15
ExecStart=/bin/bash /home/nkrueger/dev/data_log/scripts/launch_node_broadcast.sh
User=nkrueger
WorkingDirectory=/home/nkrueger/dev/data_log
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
